# com.mat.github-review-monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mat.github-review-monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mat/git/pr-notifier/github-review-monitor.sh</string>
    </array>
    
    <!-- Schedule for daily at 8:01 AM -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>1</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/tmp/github-review-monitor.out</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/github-review-monitor.err</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
        <key>HOME</key>
        <string>/Users/mat</string>
        <key>USER</key>
        <string>mat</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>